# netjack_pifmrds.service
[Unit]
Description=Netjack-PiFMrds
After=network.target


[Service]
ExecStart=/usr/local/bin/netjack_pifmrds


[Install]
WantedBy=multi-user.target
